# conf/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
pidfile=/var/run/supervisord.pid
logfile=/var/log/supervisord.log

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:xvfb]
priority=10
directory=/
command=/usr/bin/Xvfb :1 -screen 0 1024x768x16
user=tribler
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/xvfb.log
redirect_stderr=true

[program:openbox]
priority=15
directory=/home/tribler
command=/usr/bin/openbox
user=tribler
autostart=true
autorestart=true
stopsignal=QUIT
environment=DISPLAY=":1",HOME="/home/tribler",USER="tribler"
stdout_logfile=/var/log/openbox.log
redirect_stderr=true

[program:x11vnc]
priority=20
directory=/
command=x11vnc -display :1 -xkb -noxrecord -noxfixes -noxdamage -forever -shared -rfbport 5900
user=tribler
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/x11vnc.log
redirect_stderr=true

[program:novnc]
priority=25
directory=/opt/noVNC/
command=/opt/noVNC/utils/launch.sh --listen 6081 --vnc 127.0.0.1:5900
user=tribler
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/novnc.log
redirect_stderr=true
stopasgroup=true

[program:tribler]
priority=20
directory=/opt/tribler/
command=/opt/tribler/tribler.sh
user=tribler
autostart=true
autorestart=true
stopsignal=QUIT
environment=DISPLAY=":1",HOME="/home/tribler",USER="tribler"
stdout_logfile=/var/log/tribler.log
redirect_stderr=true
stopasgroup=true

[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
user=root
autostart=true
autorestart=true
priority=30
stopsignal=QUIT
;stdout_events_enabled=true
;stderr_events_enabled=true
stdout_logfile=/var/log/nginx.log
redirect_stderr=true
;stdout_logfile=/dev/stdout
;stdout_logfile_maxbytes=0
;stderr_logfile=/dev/stderr
;stderr_logfile_maxbytes=0

